# petdex-desktop.service
[Unit]
Description=PetDex Desktop Pet
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/petdex-desktop
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0
Environment=RUNTIME_DIR=%h/petdex/runtime

[Install]
WantedBy=graphical-session.target
